# deploy/systemd/homeops-backup-jsonl.timer
[Unit]
Description=Daily HomeOps JSONL Log Backup
Documentation=https://github.com/dhleach/homeops/blob/master/deploy/backup/backup-jsonl-logs.sh
Requires=homeops-backup-jsonl.service

[Timer]
# Run daily at 3 AM UTC (10 PM EDT)
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
